# window-manager.target
[Unit]
Description=X11 window manager
RefuseManualStart=yes
