# x11.target
[Unit]
Description=Graphical interface
Wants=window-manager.service
